# rootfs/etc/services.d/ha-mcp/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
exec python3 -m app.server
